# suspend.timer
[Unit]
Description=suspend timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
